# s6/etc/services.d/vnc/finish
#!/usr/bin/with-contenv sh

s6-echo "Stopping VNC server"
pkill -9 x11vnc
