# s6/xvfb/finish
#!/command/with-contenv sh

s6-echo "Stopping Xvfb"
pkill -9 xvfb-run
