# s6/fluxbox/run
#!/command/with-contenv sh

s6-echo "Starting fluxbox"
/usr/bin/fluxbox -display ${DISPLAY} >/var/log/fluxbox.log 2>&1
